# bin/callattendant.service.example
[Unit]
Description=Call Attendant
After=multi-user.target

[Service]
Type=simple
ExecStart=callattendant --config /home/pi/.callattendant/app.cfg
WorkingDirectory=/home/pi/.callattendant
Restart=on-abort

[Install]
WantedBy=multi-user.target
